# platform/nginx-nodejs/docker/config/supervisor/conf.d-sample/nodejs-index.conf
[program:nodejs-index]
command=node --watch /var/www/index.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0