# platforms/nginx-nodejs-25/docker/config/supervisor/conf.d-sample/nestjs-dev.conf
[program:nestjs-dev]
command=npx nest start --watch
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0